# seafile/seafile-memcached.container
[Container]
ContainerName=seafile-memcached
Exec=memcached -m 1024
Image=docker.io/library/memcached:1.6.27
Network=seafile.network
